# .config/systemd/user/pkglist.service
[Service]
ExecStart=/bin/make -C %h/.config/pkglist
Type=oneshot
Restart=no
