# .config/systemd/user/pkglist-watch-packages.path
[Path]
Unit=pkglist.service
PathChanged=/var/lib/pacwall/systemd/trigger

[Install]
WantedBy=default.target
